# backup/crontab
17 4 * * * /backup/backup >> /var/log/cron.log 2>&1
